# systemd/demo-app.service
[Unit]
Description=Minimal WSL Python Demo App
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/workspace/wsltest
ExecStart=/root/workspace/wsltest/venv/bin/python3 /root/workspace/wsltest/main.py --loop --interval 10
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
